# deploy/nebula_server.conf
[program:nebula_server]
directory=/nebula
command=/usr/local/bun/gunicorn -c /nebula/deploy/gunicorn.conf.py main:app
user=root
autostart=true
autorestart=true
startretries=5
redirect_stderr=true
stdout_logfile=/var/log/nebula_server/nebula_server.log
